# supervisord_base.conf
[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
pidfile=/tmp/supervisord.pid
childlogdir=/workspace/logs

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; =============================================================================
; LLM: Qwen3.6-27B INT4 AutoRound (vllm) - port 8000
;
; DEFAULT: Intel/Qwen3.6-27B-int4-AutoRound
;   - Highest quality INT4 quant (~95-97% of FP8)
;   - ~14GB weights, ~127GB KV cache headroom on H200
;   - Quantization: gptq_marlin, dtype: float16
;
; ALTERNATIVE: Qwen/Qwen3.6-27B-FP8 (highest quality, slightly slower)
;   To switch, set RunPod env var: LLM_MODEL=Qwen/Qwen3.6-27B-FP8
;   And remove --quantization gptq_marlin, change --dtype to auto
;
; V1 engine enabled, CUDA graphs active (enforce-eager removed).
; Qwen3.6-27B dense model with gated delta networks hybrid attention + MTP.
; H200 (141GB HBM3e) - INT4 ~14GB weights, ~127GB KV cache headroom.
; =============================================================================
[program:vllm]
command=vllm serve %(ENV_LLM_MODEL)s
    --quantization gptq_marlin
    --dtype float16
    --host 0.0.0.0
    --port 8000
    --gpu-memory-utilization 0.95
    --max-model-len 65536
    --kv-cache-dtype fp8
    --enable-chunked-prefill
    --max-num-batched-tokens 8192
    --enable-prefix-caching
    --language-model-only
    --reasoning-parser qwen3
    --speculative-config '{"method":"mtp","num_speculative_tokens":3}'
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
stdout_logfile=/workspace/logs/vllm.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/workspace/logs/vllm.err
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
environment=HF_HOME="/workspace/huggingface",VLLM_ATTENTION_BACKEND="FLASHINFER",VLLM_LOGGING_LEVEL="INFO",VLLM_USE_V1=1,PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True"

; =============================================================================
; TTS section is written dynamically by /start.sh based on TTS_BACKEND env var.
; See start.sh for details.
; =============================================================================
