# supervisord_base.conf
[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
pidfile=/tmp/supervisord.pid
childlogdir=/workspace/logs

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; =============================================================================
; LLM: Qwen3.5-27B-FP8 (vllm) - port 8000
; V1 engine enabled, CUDA graphs active (enforce-eager removed).
; =============================================================================
[program:vllm]
command=vllm serve Qwen/Qwen3.5-27B-FP8
    --quantization fp8
    --dtype auto
    --host 0.0.0.0
    --port 8000
    --gpu-memory-utilization 0.85
    --max-model-len 16384
    --kv-cache-dtype auto
    --enable-chunked-prefill
    --max-num-batched-tokens 8192
    --enable-prefix-caching
    --language-model-only
    --speculative-config '{"method":"qwen3_next_mtp","num_speculative_tokens":2}'
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
stdout_logfile=/workspace/logs/vllm.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/workspace/logs/vllm.err
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
environment=HF_HOME="/workspace/huggingface",VLLM_ATTENTION_BACKEND="FLASHINFER",VLLM_LOGGING_LEVEL="INFO",VLLM_USE_V1=1,PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True"

; =============================================================================
; TTS section is written dynamically by /start.sh based on TTS_BACKEND env var.
; See start.sh for details.
; =============================================================================
