# supervisord.conf
[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
pidfile=/tmp/supervisord.pid
childlogdir=/workspace/logs

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; =============================================================================
; LLM: Qwen3.5-35B-A3B-AWQ (vllm) - port 8000
; GPU memory utilization reduced to 0.85 to leave headroom for TTS process.
; V1 engine enabled, CUDA graphs active (enforce-eager removed).
; =============================================================================
[program:vllm]
command=vllm serve cyankiwi/Qwen3.5-35B-A3B-AWQ-4bit
    --quantization awq
    --quantization compressed-tensors
    --dtype auto
    --host 0.0.0.0
    --port 8000
    --gpu-memory-utilization 0.85
    --max-model-len 32768
    --kv-cache-dtype fp8
    --enable-chunked-prefill
    --max-num-batched-tokens 8192
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
stdout_logfile=/workspace/logs/vllm.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/workspace/logs/vllm.err
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
environment=HF_HOME="/workspace/huggingface",VLLM_ATTENTION_BACKEND="FLASHINFER",VLLM_LOGGING_LEVEL="INFO",VLLM_USE_V1=1,PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True"

; =============================================================================
; TTS: Qwen3-TTS-12Hz-1.7B-Base (vllm-omni) - port 8091
; Base model supports voice cloning via ref_audio/ref_text.
; stage-configs-path uses our optimized config with enforce_eager: false on both stages.
; Streaming: POST /v1/audio/speech with stream=true, response_format=pcm
;   -> raw 16-bit signed PCM at 24kHz, chunked HTTP transfer.
; Voice cloning: pass task_type=Base, ref_audio (URL or path), ref_text
; CUDA graphs are now controlled per-stage in qwen3_tts_optimized.yaml (enforce_eager: false).
; The --no-enforce-eager CLI flag is NOT what controls this for vllm-omni - the stage config is.
; =============================================================================
[program:qwen3-tts]
command=/bin/bash -c 'VLLM_LOGGING_LEVEL=INFO vllm-omni serve Qwen/Qwen3-TTS-12Hz-0.6B-Base --stage-configs-path /etc/qwen3_tts_optimized.yaml --omni --host 0.0.0.0 --port 8091 --trust-remote-code'
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
stdout_logfile=/workspace/logs/qwen3-tts.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/workspace/logs/qwen3-tts.err
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
environment=HF_HOME="/workspace/huggingface",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True"
